# deployment-srv/roles/install_redis/templates/redis.service
[Unit]
Description=Redis

[Service]
ExecStart=/opt/redis/src/redis-server /opt/redis/redis.conf

[Install]
WantedBy=multi-user.target

